# telemetry_service/telemetry_service.service
[Unit]
Description=Multi-UAV Telemetry Service
Documentation=https://github.com/your-repo/haberlesme_projesi/blob/main/README.md
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
# This path should point to where you install the executable
ExecStart=/usr/local/bin/telemetry_service
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5

# It's recommended to run the service as a dedicated non-root user.
# First, create the user: `sudo useradd --system --home-dir /var/lib/telemetry --create-home telemetry`
# Then uncomment the following lines:
# User=telemetry
# Group=telemetry

# The service will look for the config in this environment variable path first.
Environment="SERVICE_CONFIG=/etc/telemetry_service/service_config.json"

# The service will write logs to the directory specified by its internal config,
# but its stdout/stderr will be captured by the systemd journal.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telemetry_service

# Grant write access to the log directory
ReadWritePaths=/var/log/telemetry_service/

[Install]
WantedBy=multi-user.target
